# cxlinux/cx-stacks/cx_stacks/templates/systemd/gunicorn.socket.j2
[Unit]
Description=Gunicorn socket for {{ name }}

[Socket]
ListenStream=/run/gunicorn/{{ name }}.sock
SocketUser={{ user | default('www-data') }}

[Install]
WantedBy=sockets.target
